# test/supervisord.conf
[unix_http_server]
file=/Users/zhangyanpeng/go-project/demo/supervisor.sock

[supervisord]
logfile = /Users/zhangyanpeng/go-project/demo/logs
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /Users/zhangyanpeng/go-project/demo/supervisor.pid
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///Users/zhangyanpeng/go-project/demo/supervisor.sock

[program:test]
command=/Users/zhangyanpeng/go-project/demo/test/2
directory=/Users/zhangyanpeng/go-project/demo/test
autostart=true
autorestart=true
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=TERM
redirect_stderr=false
stdout_logfile=/Users/zhangyanpeng/go-project/demo/test/logs/console-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=1
stdout_capture_maxbytes=0
stdout_events_enabled=false
stderr_logfile=/Users/zhangyanpeng/go-project/demo/test/logs/console-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=1
stderr_capture_maxbytes=0
stderr_events_enabled=false
